# .docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start
autostart=true
autorestart=true

[program:php7.3-fpm]
command=service php7.3-fpm start
autostart=true
autorestart=true